# crontab
* * * * * /home/start.sh
# Mandatory blank line